# service/95-afc.rules
SUBSYSTEM=="pci", ATTR{vendor}=="0x10ee", ATTR{subsystem_device}=="0x0007", RUN+="/opt/afc-epics-ioc/service/ioc-start.sh %k"
